# etc/cron.d/pinpoint
# Pinpoint - Cron jobs
# Update lists every 6 hours to keep routing fresh
0 */6 * * * root /usr/bin/python3 /opt/pinpoint/scripts/pinpoint-update.py update >/dev/null 2>&1 || /opt/pinpoint/scripts/pinpoint-update.sh update >/dev/null 2>&1
