# etc/cron.d/pinpoint
# Pinpoint - Cron jobs
# Update lists daily at 4:00 AM
0 4 * * * root /usr/bin/python3 /opt/pinpoint/scripts/pinpoint-update.py update >/dev/null 2>&1
