# etc/init.d/pinpoint
#!/bin/sh /etc/rc.common
# Pinpoint - Selective routing service for OpenWrt (Lite Mode)

START=99
STOP=10
USE_PROCD=1

PINPOINT_DIR="/opt/pinpoint"

start_service() {
    logger -t pinpoint "Starting pinpoint service..."
    
    # 1) Start sing-box first so tun1 exists before init
    if [ -f /opt/pinpoint/data/tunnels.json ]; then
        /etc/init.d/sing-box enable >/dev/null 2>&1 || true
        /etc/init.d/sing-box start >/dev/null 2>&1 || true
        for _ in 1 2 3 4 5 6 7 8 9 10; do
            ip link show tun1 >/dev/null 2>&1 && break
            sleep 1
        done
    fi
    
    # 2) Initialize nftables and policy routing
    /opt/pinpoint/scripts/pinpoint-init.sh start
    
    # 3) Apply rules (creates dnsmasq config and loads IPs)
    /opt/pinpoint/scripts/pinpoint-apply.sh reload
    
    if [ -f /etc/dnsmasq.d/pinpoint.conf ]; then
        /etc/init.d/dnsmasq restart >/dev/null 2>&1 || true
    fi
    
    logger -t pinpoint "Pinpoint service started"
}

stop_service() {
    logger -t pinpoint "Stopping pinpoint service..."
    /opt/pinpoint/scripts/pinpoint-init.sh stop
    logger -t pinpoint "Pinpoint service stopped"
}

reload_service() {
    logger -t pinpoint "Reloading pinpoint rules..."
    /opt/pinpoint/scripts/pinpoint-apply.sh reload
}

status_service() {
    /opt/pinpoint/scripts/pinpoint-init.sh status
}
